# templates/etc.systemd.system.iqserver.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nexus IQ Service
After=network.target

[Service]
Type=simple
ExecStart=/etc/init.d/iqserver start
ExecStop=/etc/init.d/iqserver stop
User={{ iqserver_user }}
Restart=on-abort

[Install]
WantedBy=multi-user.target